# dot_config/systemd/user/podman-auto-update.timer.d/override.conf
[Timer]
# Clear the existing schedule
OnCalendar=
# Set to run every 6 hours
OnCalendar=TUE *-*-* 00:00:00
Persistent=true
